# data/supervisor.conf
[program:smsgateway]
command = python /opt/loxberry/data/plugins/smsgateway/scripts/smsout.py
directory = /opt/loxberry/data/plugins/smsgateway/scripts
user = root
environment=HOME="/root", USER="root"
